# models/pineapple-example-modules/src/main/resources/pineapple-example-modules/weblogic-wlst-009-r121-configure-javabased-nodemanager/bin/nodemanager
#!/bin/sh
#
# chkconfig: 235 91 35
# description: starts and stops the node manager
#
#
. /etc/rc.d/init.d/functions

RETVAL=0
SERVICE="nodemanager"
NM_USER=weblogic
LOCKFILE=/var/lock/subsys/${SERVICE}
HOME_DIR=/u01/app/oracle/product/fmw/wlserver_12.1
BIN_DIR=${HOME_DIR}/server/bin
LOG_DIR=/tmp/
LOG_FILE=${LOG_DIR}/nodemanager.log

# See how we were called.
case "$1" in
  start)
        # Start daemon.
        echo -n "Starting NodeManager: "
        su -p -s /bin/sh $NM_USER -c "$BIN_DIR/startNodeManager.sh  > $LOG_FILE &"
        RETVAL=$?
        [ $RETVAL -eq 0 ] && success || failure
        echo
        [ $RETVAL -eq 0 ] && touch $LOCKFILE
        ;;
  stop)
        # Stop daemons.
        echo -n "Shutting down NodeManager: "
        PID=`ps -ef | grep "${SERVICE}" | grep "${HOME_DIR}" | grep -v grep | grep -v start | cut -c 10-14`
        kill -15 $PID
        RETVAL=$?
        echo
        [ $RETVAL = 0 ] && rm -f $LOCKFILE
        ;;
  restart)
        $0 stop
        $0 start
        ;;
  status)
        PCOUNT=`ps -ef | grep "${SERVICE}" | grep "${HOME_DIR}" | grep -v grep | grep -v start | wc -l`
        if [ "$PCOUNT" -eq 0 ]
        then
                echo $"Service is not running"
        else
                echo $"Service is running in `echo $PCOUNT` processes"
        fi
        ;;
  *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
esac

exit  $RETVAL